# pick_pocket.crontab
*/5 * * * * node /FULL_PATH/pick_pocket/app.js